# monitor/config/upstart/statsd.conf
description "StatsD"

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

script
  exec su -l ubuntu -c "/home/ubuntu/PubnubTest/monitor/scripts/launch_statsd.sh"
end script
